# launchd/com.aisechub.daily-curator.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.aisechub.daily-curator</string>

  <key>ProgramArguments</key>
  <array>
    <string>/Applications/Codex.app/Contents/Resources/node</string>
    <string>/Users/taleliyahu/newsletter/bin/aisechub-daily-curator.mjs</string>
    <string>--limit</string>
    <string>5</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/taleliyahu/newsletter</string>

  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>9</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>

  <key>StandardOutPath</key>
  <string>/Users/taleliyahu/newsletter/logs/aisechub-daily-curator.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/taleliyahu/newsletter/logs/aisechub-daily-curator.err.log</string>
</dict>
</plist>

